# system/install/deconzWebhook.service
[Unit]
Description=deconzWebhook

[Service]
Type=idle
ExecStart=REPLACE/start_deconzWebhook.sh
Restart=on-failure

[Install]
WantedBy=default.target
