# system/rpi_zrct300_vlc.service
[Unit]
Description=rpi_zrct300_vlc
After=multi-user.target

[Service]
Type=idle
ExecStart=su USER REPLACE/start_controller.sh

[Install]
WantedBy=multi-user.target
