# system/rpi_zrct300_vlc.service
[Unit]
Description=rpi_zrct300_vlc

[Service]
Type=idle
ExecStart=REPLACE/start_controller.sh
Restart=on-failure

[Install]
WantedBy=default.target
