# system/streamdeck.service
[Unit]
Description=streamdeck

[Service]
Type=simple
ExecStart=REPLACE/start_streamdeck.sh
Restart=on-failure

[Install]
WantedBy=default.target
